# luci-app-tinyfm/root/etc/uci-defaults/tinyfm-symlink.sh
#!/bin/sh

cd /
[ ! -d /www/tinyfm/rootfs ] && ln -s / /www/tinyfm/rootfs
exit 0
