# luci-app-tinyfm/root/etc/uci-defaults/tinyfm-symlink.sh
#!/bin/sh

cd /
ln -s / /www/tinyfm/rootfs
exit 0
